# contrib/fuel_bootstrap/files/trusty/etc/rc.local
#!/bin/sh -e

# Apply work-around for net.ifnames
let-rename || true

# Perform fuel bootstrap configuration
fix-configs-on-startup || true


