# contrib/fuel_bootstrap/files/trusty/etc/rc.local
#!/bin/sh -e

# Perform fuel bootstrap configuration
fix-configs-on-startup || true


