# contrib/fuel_bootstrap/files/xenial/etc/rc.local
#!/bin/sh -e

# Perform fuel bootstrap configuration
/usr/bin/fix-configs-on-startup || /bin/true
