# contrib/fuel_bootstrap/files/xenial/etc/rc.local
#!/bin/sh -e

# Perform fuel bootstrap configuration
fix-configs-on-startup || true

# Unlocking nailgun-agent for discovering
rm -f /etc/nailgun-agent/nodiscover
