# contrib/ironic/bootstrap-files/lib/systemd/system/ironic-callback.service
[Unit]
After=sshd.service

[Service]
ExecStart=/usr/bin/ironic-callback.sh

[Install]
WantedBy=default.target
